# root/etc/services.d/aria2/run
#!/usr/bin/with-contenv sh
exec aria2c --enable-rpc --rpc-listen-all --rpc-secret=$ARIA2_SECRET_TOKEN
