# root/etc/services.d/aria2/run
#!/usr/bin/with-contenv sh
exec aria2c --conf-path=/etc/aria2/aria2.conf
